# RaspberryPi/resetbuttond/resetbuttond.service.in
[Unit]
Description=Reset button handler

[Service]
Type=simple
ExecStart=@CMAKE_INSTALL_PREFIX@/sbin/resetbuttond 18 24

[Install]
WantedBy=default.target
