# RaspberryPi/xbee-tools/xbeemux.service.in
[Unit]
Description=XBee multiplexer service

[Service]
Type=simple
ExecStart=@CMAKE_INSTALL_PREFIX@/sbin/xbeemux /dev/ttyAMA0

[Install]
WantedBy=default.target
